# deploy/stremio-addarr.service.example
[Unit]
Description=stremio-addarr – Stremio Arr Status & Add add-on
Documentation=https://github.com/cbkii/stremio-addarr
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/stremio-addarr
EnvironmentFile=/opt/stremio-addarr/.env
ExecStart=/usr/bin/node /opt/stremio-addarr/dist/src/index.js
Restart=on-failure
RestartSec=5

# REQUIRED: replace with the service account/group that exist on your host.
# Do not leave placeholders in production.
User=YOUR_SVC_USER
Group=YOUR_SVC_GROUP

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/stremio-addarr
PrivateTmp=true

[Install]
WantedBy=multi-user.target
